# crontab
#run python script every minutes
*/10 * * * * python3 /app/aliyun_ddns.py > /proc/1/fd/1 2>/proc/1/fd/2
